# roles/base/files/etc/systemd/system/routing-rulesd.service
[Unit]
Description="Initialize routing rules"

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/routing-rulesd

[Install]
WantedBy=network-pre.target
